# resources/etc/supervisor/conf.d/bowhead.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server
startretries=3
stopwaitsecs=10
autostart=true

